# packages/shaping-scripts/shaper@.service
[Unit]
Description=Traffic shaping on interface %I
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/system/shaper@.sh %I up
ExecStop=/usr/lib/systemd/system/shaper@.sh %I down
RemainAfterExit=yes

[Install]
WantedBy=network-online.target
