# packages/mbsync-service/mbsync@.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=mbsync@%i.service

[Install]
WantedBy=timers.target
