# packages/snapper-backup/snapper-backup.service
[Unit]
Description=Backup of root snapshots
After=snapper-timeline.service
After=snapper-cleanup.service

[Service]
Type=simple
ExecStart=/usr/lib/systemd/system/snapper-backup.py /.snapshots /tank/Backups

[Install]
WantedBy=snapper-timeline.service
WantedBy=snapper-cleanup.service
